# local-luigi-workflow/files/crontab/monitor-job.crontab
#!/bin/bash
*/1 * * * * /bin/ps aux --sort -rss >> memory-monitor.log 2>&1

