# local-luigi-workflow/files/crontab/monitor-job.crontab
*/5 * * * * ps -C python -o pid,%cpu,%mem,stime,etime,command >> memory-monitor.log
